# services/weston/run
#!/command/execlineb -P
with-contenv weston --xwayland -B headless